# image-files/var/spool/cron/crontabs/root
PATH="/usr/local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
*/5 * * * * /usr/local/bin/python3 /var/work/ip.py >> /var/log/cron.log 2>&1
